# pydhcp/pydhcpd.service
[Unit]
Description=pydhcpd - Python DHCP Daemon
Documentation=https://github.com/maravento/vault/tree/master/pydhcp
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=10

[Service]
Type=simple
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONDONTWRITEBYTECODE=1"
ExecStart=/usr/bin/python3 /etc/pydhcp/pydhcpd.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
User=pydhcpd
Group=pydhcpd
AmbientCapabilities=CAP_NET_RAW CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/etc/pydhcp /var/log

[Install]
WantedBy=multi-user.target
