# gateproxy/conf/pack/ttyd.service
[Unit]
Description=ttyd - Web Terminal
After=network.target

[Service]
Type=simple
User=your_user
WorkingDirectory=/home/your_user
ExecStart=/usr/bin/ttyd -i lo -p 7681 -W bash
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
